# hypr/scripts/awww-wallpaper.service
[Unit]
Description=Change wallpaper

[Service]
Type=oneshot
ExecStart=%h/.config/hypr/scripts/change-wallpaper.sh
